# ansible/roles/darkedges.kubernetes/templates/etc/profile.d/kubernetes.sh.j2
#!/bin/bash

if [ ! -d "$HOME/.kube" ]; then
    mkdir -p $HOME/.kube
    sudo cp -i /etc/kubernetes/admin.conf $HOME/.kube/config
    sudo chown $(id -u):$(id -g) $HOME/.kube/config
fi    